# tempergold-temp-recorder.timer
[Unit]
Description=TemperGold Schedule

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target

